# app/crontab
* * * * * /usr/local/bin/python /app/main.py >> /app/app.log 2>&1
* * * * * /app/downlatestimg.sh >> /app/app.log 2>&1